# etc/systemd/system/dnf-automatic.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/dnf clean metadata
ExecStopPost=/usr/local/bin/divvun-apy-extract-all-zcheck
ExecStopPost=/usr/bin/systemctl restart apy
