# etc/systemd/system/dnf-automatic.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/dnf clean metadata
ExecStopPost=/usr/bin/unzip -o -u -d /usr/share/voikko/4/se /usr/share/voikko/4/se.zcheck
ExecStopPost=/usr/bin/systemctl restart apy
