# etc/systemd/system/dnf-automatic.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/dnf clean metadata
ExecStopPost=/usr/bin/systemctl restart apy
