# etc/systemd/system/certbot.timer
[Unit]
Description=Check if letsencrypt cert needs renewal

[Timer]
OnCalendar=04,21:21

[Install]
WantedBy=timers.target

