# etc/systemd/system/restart-apy.service
[Unit]
Description=Restart APY server

[Service]
Type=simple
Restart=no
ExecStart=/usr/bin/systemctl restart apy.service

