# etc/systemd/system/svn-up-apy.service
[Unit]
Description=SVN up apy user's checkouts

[Service]
Type=simple
Restart=no
ExecStart=/usr/bin/svn up /home/apy/GTHOME

