# etc/systemd/system/svn-up-apy.timer
[Unit]
Description=SVN up apy user's checkouts every so often

[Timer]
OnCalendar=02,08,12,14,16,18:42

[Install]
WantedBy=timers.target

